# examples/nova-vpn-bot.service
[Unit]
Description=Nova VPN Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nova-vpn
ExecStart=/usr/bin/python3 bot/main.py
EnvironmentFile=/opt/nova-vpn/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
